# xkbset-bouncekeys.service
[Unit]
Description=xkbset bouncekeys

[Service]
ExecStart=/home/ncaq/.rc/xkbset-bouncekeys
Type=oneshot

[Install]
WantedBy=graphical.target
